# web/php7.2_ice3.7.3/crontab/cron
* * * * * root echo hi >> /root/log.log
